# oracle9-myapp/rootfs/etc/systemd/system/myapp.service
[Unit]
Description='my app'

[Service]
ExecStart=/opt/myapp/app
Type=oneshot

[Install]
WantedBy=multi-user.target
